# root/etc/s6-overlay/s6-rc.d/svc-mariadb/finish
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

echo "Caught SIGTERM signal!"

kill -TERM "$(cat /var/run/mysqld/mysqld.pid)" 2>/dev/null
