# services/mysql/run
#!/bin/bash
/sbin/setuser abc /usr/bin/mysqld_safe --skip-syslog
