# root/etc/services.d/mariadb/run
#!/usr/bin/with-contenv bash

exec \
	/usr/bin/mysqld --user=abc --datadir="$DATADIR"
