# root/etc/services.d/mariadb/run
#!/usr/bin/with-contenv bash
 exec \
	s6-setuidgid abc /usr/bin/mysqld_safe
